# .config/systemd/user/mpris-proxy.service
[Unit]
Description=MPRIS proxy

[Service]
ExecStart=/usr/bin/mpris-proxy

[Install]
WantedBy=default.target
